# services/voice.service
[Unit]
Description=PiCar Voice Assistant (Hey Jarvis)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/picar-brain
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /home/pi/picar-brain/voice_assistant.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
